# src/etc/cron.d/tsidp
*/5 * * * * /usr/local/php/unraid-tsidp/check-tsidp.php >> /var/log/tsidp-check.log 2>&1